# deploy/travel-planner-web.service
# systemd 服务单元 —— WanderWarm Web 前端 (FastAPI + SSE, 端口 8000)
#
# 安装:
#   cp deploy/travel-planner-web.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now travel-planner-web
#
# 与 travel-planner.service (Streamlit, 8501) 互不影响, 可同时运行。
# 记得在阿里云防火墙放行 TCP 8000。

[Unit]
Description=Multi-Agent Travel Planner (WanderWarm Web)
After=network.target

[Service]
WorkingDirectory=/root/multi-agent-travel-planner/python
ExecStart=/root/multi-agent-travel-planner/python/.venv/bin/uvicorn api.app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
